# debian/abraflexi-digest.cron.weekly
#!/bin/sh
/usr/bin/flexibee-weekdigest

