# debian/abraflexi-digest.cron.monthly
#!/bin/sh
/usr/bin/flexibee-monthdigest

