# debian/abraflexi-digest.cron.d
#!/bin/sh
55 23 28-31 * * [[ "$(date --date=tomorrow +\%d)" == "01" ]] &&	/usr/bin/flexibee-yeardigest

